# roles/selenium-hub/templates/selenium-hub.service.j2
[Unit]
Description=Selenium Hub

[Service]
Type=simple
ExecStart=/usr/bin/java {{ java_opts }} -jar /opt/selenium/{{ selenium_jar }} -role hub -hubConfig /opt/selenium/config.json